# system/docker-network-create.service
[Unit]
Description=Docker Network Creation
Requires=docker.service
After=docker.service

[Service]
ExecStart=/usr/bin/docker network create nats